# linux_udev_rules/30-pcscd.rules
ACTION=="add", SUBSYSTEM=="usb",ENV{DEVTYPE}=="usb_device", ATTRS{product}=="ACR122U PICC Interface", RUN+="/usr/bin/reset-demo"
