# linux_udev_rules/30-pcscd.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{product}=="ACR122U PICC Interface", RUN+="/usr/bin/reset-demo"
