# network/config-network.service
[Unit]
Description=Configure the network of the Valeo charger 
After=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/config-network

[Install]
WantedBy=multi-user.target
